# nodes/x86/etc/supervisor.d/agent.ini

[program:openvpn]
priority=10
directory=/
command=/usr/sbin/openvpn --config /etc/openvpn/docker.ovpn
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/openvpn.log
redirect_stderr=true

[program:socat]
priority=15
directory=/
command=/usr/bin/socat -T 120 tcp4-listen:23750,fork unix-connect:/var/run/docker.sock
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/socat.log
redirect_stderr=true

[program:swarm-agent]
priority=20
directory=/
command=/bin/swarm.sh consul://172.20.0.1:8500 
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/swarm-agent.log
redirect_stderr=true

[program:ip-checker]                   
priority=25                            
directory=/         
command=/bin/checkip.sh
user=root              
autostart=true         
autorestart=true       
stdout_logfile=/var/log/ip-checker.log                          
redirect_stderr=true                   
                    
